# Configs/V2Ray/cloudflared_panel.service
[Unit]
Description=Cloudflared FastAPI Panel
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/runner
ExecStart=/root/.local/bin/uvicorn cfagent:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
