# V2Ray/6h.timer
[Unit]
Description=Timer for 6h.service

[Timer]
OnBootSec=5h59min
AccuracySec=1s
Persistent=false

[Install]
WantedBy=timers.target
